# files/root/etc/cron.d/check-cert
0 4 * * * sh /check_cert.sh >> /var/log/cron.log 2>&1